# startup/newbsd-init.in
#!/bin/sh
#
# Copyright (c) 2016 Nagios(R) Core(TM) Development Team
#
# PROVIDE: nrpe
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

sig_reload=HUP
: ${nrpe@bsdenable@:="NO"}
: ${nrpe_configfile:="@pkgsysconfdir@/nrpe.cfg"}
: ${nrpe_program:="@sbindir@/nrpe"}
: ${pidfile:="@piddir@/nrpe.pid"}
command_args="-c ${nrpe_configfile} -d"
extra_commands=reload
name=nrpe
rcvar=nrpe@bsdenable@
load_rc_config "${name}"
required_files="${nrpe_configfile}"
sig_reload=HUP

start_precmd=nrpe_prestart
stop_precmd=find_pidfile

nrpe_prestart()
{
	[ -n "$nrpe_pidfile" ] &&
		warn "No longer necessary to set nrpe_pidfile in rc.conf[.local]"

	install -d -o @nrpe_user@ ${pidfile%/*}
}

run_rc_command "$1"
