# startup/openrc-init.in
#!/sbin/openrc-run
#
# Copyright (c) 2017 Nagios(R) Core(TM) Development Team
#

# Supply a default value for NRPE_CFG in case the corresponding
# conf.d file is not installed.
: ${NRPE_CFG:="@sysconfdir@/nrpe.cfg"}

command="@sbindir@/nrpe"
command_args="--config=${NRPE_CFG} ${NRPE_OPTS}"
command_args_background="--daemon"
description="Nagios Remote Plugin Executor (NRPE) daemon"
extra_started_commands="reload"
pidfile="@piddir@/${RC_SVCNAME}.pid"

reload() {
    ebegin "Reloading ${RC_SVCNAME}"
    start-stop-daemon --signal HUP --pidfile "${pidfile}"
    eend $?
}
