# startup/solaris-PKG_NAME.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<!--
	This is a manifest file for the nrpe service. Import it using
		svccfg import <path/to/this/file>
	and enable it using
		svcadm enable nrpe
-->

<service_bundle type='manifest' name='NGOS:nrpe'>

	<service
		name='network/nagios/nrpe'
		version='1'
		type='service'>

		<single_instance />

		<dependency
			name='fs-local'
			grouping='require_all'
			restart_on='none'
			type='service'>
				<service_fmri value='svc:/system/filesystem/local' />
		</dependency>

		<dependency
			name='autofs'
			grouping='optional_all'
			restart_on='none'
			type='service'>
			<service_fmri value='svc:/system/filesystem/autofs' />
		</dependency>

		<dependency name='net-loopback'
			grouping='require_all'
			restart_on='none'
			type='service'>
			<service_fmri value='svc:/network/loopback' />
		</dependency>

		<dependency name='net-physical'
			grouping='require_all'
			restart_on='none'
			type='service'>
			<service_fmri value='svc:/network/physical' />
		</dependency>

		<dependency name='cryptosvc'
			grouping='require_all'
			restart_on='none'
			type='service'>
			<service_fmri value='svc:/system/cryptosvc' />
		</dependency>

		<dependency name='utmp'
			grouping='require_all'
			restart_on='none'
			type='service'>
			<service_fmri value='svc:/system/utmp' />
		</dependency>

		<dependency name='config_data'
			grouping='require_all'
			restart_on='restart'
			type='path'>
			<service_fmri
			    value='file://localhost/pkg/etc/nrpe.cfg' />
		</dependency>

		<dependency
			name='system-log'
			grouping='optional_all'
			restart_on='none'
			type='service'>
			<service_fmri value='svc:/system/system-log' />
		</dependency>

		<dependent
			name='nrpe_multi-user-server'
			grouping='optional_all'
			restart_on='none'>
				<service_fmri
				    value='svc:/milestone/multi-user-server' />
		</dependent>

		<method_context>
			<method_credential user='nagios' group='nagios'/>
			<method_environment>
				<envvar name='BASEDIR' value='/export/home/nagios'/>
			</method_environment>
		</method_context>

		<exec_method
			type='method'
			name='start'
			exec='/usr/sfw/sbin/tcpd -c /etc/nagios/nrpe.cfg -i'
			arg0='${exec_prefix}/sbin/nrpe'
			timeout_seconds='60'/>

		<exec_method
			type='method'
			name='stop'
			exec=':kill'
			timeout_seconds='60'/>

		<exec_method
			type='method'
			name='refresh'
			exec=':hup'
			timeout_seconds='60'/>

		<property_group name='startd' type='framework'>
			<propval name='ignore_error'
				type='astring'
				value='core,signal' />
		</property_group>

		<property_group name='general' type='framework'>
			<propval name='enabled'
				type='boolean'
				value='false'/>
			<propval name='action_authorization'
				type='astring'
				value='solaris.smf.manage.nrpe'/>
			<propval name='value_authorization'
				type='astring'
				value='solaris.smf.manage.nrpe'/>
		</property_group>

		<stability value='Unstable' />

		<create_default_instance name='default' enabled='false' />

		<template>
			<common_name>
				<loctext xml:lang="C">NRPE daemon</loctext>
			</common_name>
			<description>
				<loctext xml:lang="C">
					Nagios Remote Plugin Executor daemon
				</loctext>
			</description>
			<documentation>
				<doc_link name='nagios.org' uri='http://www.nagios.org' />
			</documentation>
		</template>

	</service>

</service_bundle>
