# package/solaris/pkg/nrpe.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--

	All source code, binaries, documentation, and information contained
	in this distribution are provided AS IS with NO WARRANTY OF ANY KIND,
	INCLUDING THE WARRANTY OF DESIGN, MERCHANTABILITY, AND FITNESS FOR
	A PARTICULAR PURPOSE.

	Nagios and the Nagios logo are registered trademarks of Nagios Enterprises.
	All other trademarks, servicemarks, registered trademarks, and 
	registered servicemarks are the property of their respective owner(s).

-->

<service_bundle type='manifest' name='NGOS:nrpe'>

<service
	name='application/nagios/nrpe'
	type='service'
	version='1'>

	<create_default_instance enabled='false' />

	<single_instance />

	<dependency name='fs-local'
		grouping='require_all'
		restart_on='none'
		type='service'>
		<service_fmri
			value='svc:/system/filesystem/local' />
	</dependency>

	<dependency name='fs-autofs'
		grouping='optional_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/system/filesystem/autofs' />
	</dependency>

	<dependency name='net-loopback'
		grouping='require_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/network/loopback' />
	</dependency>

	<dependency name='net-physical'
		grouping='require_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/network/physical' />
	</dependency>

	<dependency name='cryptosvc'
		grouping='require_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/system/cryptosvc' />
	</dependency>

	<dependency name='utmp'
		grouping='require_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/system/utmp' />
	</dependency>

	<dependency name='config_data'
		grouping='require_all'
		restart_on='restart'
		type='path'>
		<service_fmri
		    value='file://localhost/etc/nagios/nrpe.cfg' />
	</dependency>

	<dependent
		name='nrpe_multi-user-server'
		grouping='optional_all'
		restart_on='none'>
			<service_fmri
			    value='svc:/milestone/multi-user-server' />
	</dependent>

	<exec_method
		type='method'
		name='start'
		exec='/lib/svc/method/nrpe start'
		timeout_seconds='60'/>

	<exec_method
		type='method'
		name='stop'
		exec='/lib/svc/method/nrpe stop'
		timeout_seconds='60' />

	<exec_method
		type='method'
		name='refresh'
		exec='/lib/svc/method/nrpe restart'
		timeout_seconds='60' />

	<property_group name='startd'
		type='framework'>
		<!-- sub-process core dumps shouldn't restart session -->
		<propval name='ignore_error'
		    type='astring' value='core,signal' />
	</property_group>

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
			NRPE daemon
			</loctext>
		</common_name>
		<description>
			<loctext xml:lang='C'>
			Nagios Remote Plugin Executor Daemon
			</loctext>
		</description>
		<!--documentation-->
			<!--manpage title='sshd' section='1M' manpath='/usr/share/man' /-->
		<!--/documentation-->
	</template>

</service>

</service_bundle>
